# crontab
* * * * * /usr/local/bin/php /var/www/html/app/artisan schedule:run >> /var/log/cron.log 2>&1
